# rootfs/etc/s6/sv/rc4/up
#!/bin/execlineb -P

s6-rc-run-init_d 4
